# docker/bytebound.service
[Unit]
Description=ByteBound Server
After=docker.service
Requires=docker.service

[Service]
User=ubuntu
WorkingDirectory=/app

Environment=COMPOSE_FILE=/app/docker/docker-compose.yml
Environment=SERVICE=bytebound

# Graceful stop: systemd sends SIGTERM, which triggers Hyper's graceful shutdown
ExecStart=/usr/bin/docker compose up -d server
ExecReload=/app/update.sh
ExecStop=/usr/bin/docker compose down

Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
